# PingRelay.supervisor.conf
[program:PingRelay]
directory=/replace/me/
command = /usr/bin/python /replace/me/PingRelay.py
autostart=true
autorestart=true
startretries=3
